# cron/crontab
#crontab
0 0 * * * web2py *cron/send_watch_email
